# ansible/roles/backup/files/backup.timer
[Unit]
Description=Run backup every day at night

[Timer]
Unit=backup.service
OnCalendar=*-*-* 03:00:00

[Install]
WantedBy=timers.target